# ansible/roles/dev-vault-init/files/auto-unseal-vault.service
[Unit]
Description="Unseal the HashiCorp Vault"
After=vault.service
ConditionFileNotEmpty=/etc/vault.d/vault.hcl

[Service]
User=vagrant
Group=vagrant
ExecStart=/home/vagrant/auto-unseal-vault.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
